# docker/backend/files/etc/crontab
0 0 * * * www-data echo "$(date --rfc-3339=seconds) - cron job" >> /var/www/app/var/log/crontab.log

# Make sure the crontab file ends on an empty line!
